# manifest/thgttg.com/golgafrinchans/etc/systemd/system/lava.service
[Unit]
Description=lavad
After=network-online.target
ConditionPathExists=/var/lib/lava/.lava/config/genesis.json
ConditionPathExists=/var/lib/lava/.local/bin/update-binary-lava.sh

[Service]
User=lava
Group=lava
ExecStartPre=/var/lib/lava/.local/bin/update-binary-lava.sh
ExecStart=/var/lib/lava/.local/bin/lavad start \
    --abci socket \
    --address tcp://0.0.0.0:26658 \
    --api.address tcp://127.0.0.1:1317 \
    --api.enable \
    --api.max-open-connections 1000 \
    --api.rpc-max-body-bytes 1000000 \
    --api.rpc-read-timeout 10 \
    --api.rpc-write-timeout 10 \
    --api.swagger \
    --block_sync \
    --db_backend goleveldb \
    --db_dir data \
    --genesis_hash f7a0c7d2587d2bf640570309137c905eac834f0aba99f90b4c10f45ef8334583 \
    --grpc.enable \
    --grpc.address 127.0.0.1:9090 \
    --grpc-web.enable \
    --grpc-web.address 127.0.0.1:9091 \
    --home /var/lib/lava/.lava \
    --inter-block-cache \
    --moniker 'cherno more' \
    --p2p.external-address golgafrinchans.thgttg.com:26656 \
    --p2p.laddr tcp://0.0.0.0:26656 \
    --p2p.persistent_peers 0e0e01f932a124c45f7f8600e38dba445b5f5dc4@65.108.226.183:19956,a8daccaa3723c327d5ddc3c1ab037f0a87ed457d@194.163.128.12:19956,68c2a08299018cdd6455e24114f4d5b3349ca896@65.108.95.13:26656,0c56139d2f456a6f30ddf5f07292e49de7c7dbed@207.180.198.10:19956,0e203a799d85b4dcd48556f3425a280f969bcf39@65.108.140.97:21656,2701fea28d683f58ce4a4b863653b05bac831e22@65.21.88.86:56656,42eb3b0a36574b1c31cdc2d69b3af49d0713cef7@65.109.52.156:6656,7fe56f54d7468dfe1d23d986d2034b2f3a0eaa3e@88.198.39.169:56656,d46adbceae5f7c4caa90bcb3892053907d9c2f18@185.250.36.72:10656,c19965fe8a1ea3391d61d09cf589bca0781d29fd@162.19.217.52:26656,eb849730cb677b8817e5aeb6948ab01c0bb077ea@176.9.157.48:27556,d40e03d235c7004639164745f4f616827f824417@134.209.247.172:19956,eb7832932626c1c636d16e0beb49e0e4498fbd5e@65.108.231.124:20656,8f7d6a584458a6af1e344ba14b8f6b54ed080a54@188.40.62.57:26656,3f7d9f4807a6998466502779427eb051d123e81a@144.91.97.90:14456,34de91a925a17fe09b62308f807e037da5380d4b@109.199.107.213:14456,be7ebba939854c314bbcef7cc651786113b8c10a@65.108.123.237:20256,5c015859867d316d7a9d41138e6e73a78f72a78c@89.117.52.186:24656,53ff7a2726de241c87e2260ad715dadc75977974@65.108.227.112:13656,d5a5d71b327ee5697226e48404fddfd2afd1a8c2@85.239.236.2:26656,df05175604b7b63d7913ac3a999c4a3b97d75a42@146.59.68.229:26656,cffa9095798f47d3233f44ec1dd002914fe38d34@162.55.184.147:26656,8b6c99e650423fe45cb3e479b9304597a421fcec@178.63.53.161:14556,a231cb05e19ded349822df3104e6262a462244d9@173.249.52.11:19956,76375db6e75706c0de3d5918fbdb07b0ae4aba10@65.108.72.253:20656,c76b91154e9bc806eb05e945acf9bf8eb4f0011a@14.226.231.57:29656,8db3d51f1c76d2417efed26ec864d3f817c188be@162.55.234.190:14656,32d0eaa31ab8f9c2779ce9272b7a68f3d15a8e6e@213.239.207.175:40636,e44268c4246c0aec339d9e518170bb91da188433@65.109.24.78:17656,f9f4b7782030882348a465b5aa8b8fbaa83a10e1@31.220.97.253:26656,304745328e64c2ba3592cddf5b1ab3c831d92802@65.21.194.43:56656,b25f6e28372996609ade4d5e44f9018c8a8e9327@195.3.220.48:26656,fd4b8c638c4e78c81a5e814ed9deea738e5da98e@84.247.179.124:26656,40046fe63bdaa9efde27707b0d3de0bf84fedf80@86.111.48.158:26656,fcbab8252cb066966b1960a1aeb76ef7680cfd8e@65.21.32.13:26656,2aafae9e8bc330adc02ea78615fca5ac15bb8846@89.116.26.223:19956,7ee8434b107b32f764f42bcc82b6fbfc25011d8f@65.108.75.107:28656,a0eab10e3807a29ded85db18b89f73a5c978ad74@65.108.75.174:33656,d60f8d9f8e5916c3ad7d7dd2d251f274b4efc310@193.122.137.30:26656,a65de5f01394199366c182a18d718c9e3ef7f981@159.148.146.132:26656,f15a58df94b78ac1b90845d24dd06a61e20e7f00@95.111.237.91:14456,b2d8a8df9295d88717fbbfd1a79c0d6ba04e1081@88.99.254.113:14556,1b9fc4028b176b42d803c6962e39cd99b69c626a@109.199.108.83:14456,bf936e472e10bb373f584cea3090ca47ef72ba68@65.108.100.31:14456,8c71c05fd08c15ed98b6e20c197eeb9e9a42fd38@65.108.236.147:20256,ba0da985b21511c22f565bf9e7d2fbdc19fd9d92@65.109.28.165:4000,e318c84e9e27424ca768aa813d3f83213aa0ecef@65.109.69.90:2000,15fa03b72584b6a84cc0538dabba7764bfc22ced@109.199.101.63:26656,653e992118dd81a620694071c1b04efd07d33714@173.212.234.32:19956,e8e2bf51a1608ffbd4abd76f8d0472cce88e621c@65.108.79.241:36656,212039c20cbe0a2bf2533323d61775b4796ecc75@178.18.252.105:24656,8c7ac1478c112dddde0475e5a8a2a3b89e5861c1@148.251.88.145:10356,f6e40cb5f0f597f91133d1f18ebd833bd699d213@65.108.200.61:11134,7f4ab572fe7ce83e82f5dfeb5f09c39c99e1ea3a@158.220.97.152:24656,0d6983bcd192c0b4a0f61e6d849c152704e2f017@91.107.148.5:26656,3693ea5a8a9c0590440a7d6c9a98a022ce3b2455@65.109.92.79:20656,97668e8113f878373e095be7475e3cb069b02072@142.132.214.131:60656,38bbe9ad8c352d8b349001861efc4438c59d29a2@161.97.104.0:14456,97982b178ba110a38ef4c6ce9a6221dd55ec6b20@212.132.77.123:14456,290d562781dd46eaa66eb88f87b60e888550ae44@65.109.92.148:61356,8c0a204b7640bafc006cc5eb8e1c1dee9eac9292@38.242.240.215:14456,845271ef64d07bf2801b7fa357c5fe508755f57d@136.243.69.100:55856,7d71b5967e710b38d22b4f5f28738f9fc999137f@89.117.57.6:26656,e23038803b0e9d0c5cd6e79049a140d270d4fa8e@164.68.105.141:24656,a731309fe879cec15b734b56c822366384dd2589@5.9.82.137:14456,93a8d717fbd8c714bc9ccfd9c3d69fde1c6f09fd@194.163.174.231:26686,a6e08a147a9c6ce79b79d4155830adaf04d8f2e9@84.247.179.121:26656,e1e18ca56ae86e9939a28c80de6cb0fc7374f000@95.216.2.219:34656,8088940c613460fec805f594bfa263f14ffe7f3a@84.247.179.118:26656,20048b58467f63eb329cd3c62858895949521253@20.117.190.146:19956,34f3bc80e2a636147984d8bbfed7e1e29e5b022d@173.212.215.48:14456,6da596d45ae87babf92da6fc939139a3d29e6a46@31.220.98.1:26656,913656c2a2e5a8446070a6461b0a5a1786dee328@213.133.100.172:27262,07c557b393b235a7b004a4a32831e54092dc24a0@91.107.147.250:26656,552815619724bfa8bc1491786e87e71c213b6615@84.247.181.72:14456,d1730b774b7c1d52dd9f6ae874a56de958aa147e@139.45.205.60:23656,a68b367ef958c1cb42bb86c9a9eba1e4a65bc603@95.217.207.236:29456,9ec2a1e83cc8e23e774cc1b649c9450b4900c45f@134.255.244.236:26656,df763dd1fceb4e8b9118ff85e1199d4faa7ea2b5@45.13.226.108:26656,aadb410edcde075bbb449c25a9125c0b992e76ad@31.220.76.133:26656,4d259ac61e88d59e61c454a13646af1b6fab722f@173.212.240.222:14456,c067df56a578c55630515af5d747b8c2a09f272b@116.202.85.52:31656,97b3648ef143d537e6aee3b11f054a0e6b6be691@57.128.44.52:26656,c40a7bc3c7aee0428273c0bfa75fcb14bf0f44c4@65.109.90.171:30656,6d3038752141e24012e8936560a12aed12efe76b@173.249.55.113:19956,45a9e68e998d1edf81105e1be9398fdff3309f6a@65.21.109.52:26656,d92758a7624a877e22332db6dc4078b3e454ff3a@128.199.29.140:20656,4bb09193f233a7f34707cedb46632677aebeea07@161.97.138.167:14456,0ec7ad403c4caa1d3cd3f8226f456cc568eeec90@65.108.124.121:36656,6ba3b6ec03839afffa64c83e18ff80a681f4968d@65.108.194.40:21756,93b9b060b2ca055faa5f76b6b435a621303217c0@43.134.211.5:26656,27a9a6c1677d0c52bd97e0d15cd837c0f86fefc7@173.212.205.21:24656,c21cc08ff77f4fcdcf4e9d479c3a3120c00d62c7@176.9.9.174:26656,e462465a3aaed13c953255281ecff2414f422415@37.60.248.227:26656,2e4b3f6738d303695c8160b7e7cefa348cb0fa08@65.21.96.69:26956,a724b94c593241890022e204e0065d8baa67168c@116.202.227.117:14456,8264373c282bc10ed54c4360664b291239aa9e6b@158.220.105.187:24656,02fc96bda8082ca39d9db8ff0b0179085f2c6125@161.35.232.174:19956,0eb3d4bcde88b901f8a6cd4d7cd3d0346aa12161@31.220.90.176:24656,2b8855a875328e0ba06de21cbf149a63127b5ca4@207.180.218.2:19956,ac838712743edce896cda12d10e4df48d0f8f8c8@65.21.73.233:60656,bb1df8347733aedee311634727b82bbefa157523@185.137.122.90:26656,dc619645f39a299cba508b77b496ee64c2b0e250@194.242.56.204:36656,07ea0ee760c3aa43f12e265863e09e2e66364cd3@46.250.224.190:26656,f1bb78a30c9381bed392fda141a5c1f6fa4d25e6@144.76.114.49:36656,07277038190e9eb8855a49b1a13d742d18d9bea5@65.108.41.172:26656,124ab5bbb7ebf06163ef203e375e5fe015223011@217.76.51.55:14456,a131f3823dd478256606ae29b68f9c1b224abfc6@157.245.145.27:19956,ab1d80119646d89a1189c6a8604fcd2a0ffc890c@46.250.224.187:26656,43e4589ade6433329266096517faa196e0f9624d@84.247.179.126:26656,acd887559eed474121433fc2329df63f2e7ce3b3@144.91.95.185:14456,9408220cc93a84e11ba04eb19109f27b00fb4a92@65.109.19.235:20256,916e6421ea9facd9a8cfb4760eb5cdfdb56633ca@161.97.122.145:24656,e37f0dde67d2aee0a0abe76acdd2ed03bc33d4bd@109.199.107.28:14456,6c5a242470bcc1212280275ec2d709361f7a0f0c@95.216.246.78:30656,7c7edaf583eb78ce719390bbc2f838da870dd8d1@164.68.114.65:24656,41608a7be7062ead40d0606d653daed41b678eef@178.18.253.82:36656,260450340c3e6dc60f400aba846b072a2cdc16a7@173.212.250.15:20656,c92e3f253c519ce24ddd07b8e3e156bfba16b3fb@45.84.0.231:14456,bbf8ef70a32c3248a30ab10b2bff399e73c6e03c@65.21.198.100:20456,2848542ec4a20ead576ff392cf67cb7bcdd84e29@3.0.89.44:19956,f9e0513555f6362bca577aee5722f71aa372036e@184.174.35.127:36656,f655ca044eaae6955046a97bffe7136f47c73f81@38.242.241.3:24656,4d845e9656ba53581884d70fe03366c77a61097a@165.140.71.129:26656,17dabfbef1e937811bf45c0729285e5e9bf88f5c@65.108.71.136:19956,1378fe595e20bed1315e1a17b580663e79f42814@173.249.55.238:19956,42f126f001121df0589fddb64f94863c976c3bfc@83.171.249.234:30656,39db6c21168c912f606560be5a904c2e20ad842c@62.171.159.226:24656,767066580d2bf5e3fb2c7ce420021b8f810d8bae@109.199.107.99:14456,581304c0472b33768e958bf0c430aacdf75ba05b@89.117.21.186:26656,68b24abc72a7e1d80c6f199aab6a4686d413b8d0@109.199.126.161:19956,51789af428293359348020d298143946c9a0492e@5.161.132.110:26656,c927dd1f6015018a9fc7517b3417b16986e4815d@89.117.21.184:26656,80a8de391e4c253c7f79dec9d5ac0bdf4812e8f7@149.102.144.37:26656,bc31ef783b7a98367df14dc09ba4b7a0adaf0eca@31.220.98.28:26656,8b6ad3fe466038c01e6668d410c0dacecafa29a2@109.199.103.240:19956,5e17e1765a40772670a405e7f1004d8e6f0834ed@176.103.222.17:26656,aa00bddd1a947f233ab9542bc85b58896683f1bd@35.189.212.113:26656,283b68478d5ec26d7a47b43ced70e36926ea27d9@173.249.55.234:19956,57efd7bf613b27d4dfcf831fd09b0ada9c6c19c1@89.117.51.129:24656 \
    --p2p.pex \
    --p2p.seeds 3a445bfdbe2d0c8ee82461633aa3af31bc2b4dc0@testnet2-seed-node.lavanet.xyz:26656,e593c7a9ca61f5616119d6beb5bd8ef5dd28d62d@testnet2-seed-node2.lavanet.xyz:26656,eb7832932626c1c636d16e0beb49e0e4498fbd5e@lava-testnet-seed.itrocket.net:20656 \
    --pruning default \
    --rpc.laddr tcp://127.0.0.1:26657 \
    --transport socket \
    --with-tendermint \
    --log_format plain \
    --log_level info
Restart=always
RestartSec=180
LimitNOFILE=65536
LimitNPROC=65536

[Install]
WantedBy=multi-user.target
