# daemon/sync-node-state-nebra.service
[Unit]
Description=sync node state nebra
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment=HOME=/root
Environment=intent_list=nebra.systems
ExecStartPre=curl -sLo /usr/local/bin/sync-node-state-nebra.sh https://raw.githubusercontent.com/grenade/rubberneck/main/daemon/sync-node-state.sh
ExecStartPre=chmod +x /usr/local/bin/sync-node-state-nebra.sh
ExecStart=/usr/local/bin/sync-node-state-nebra.sh
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
