# static/etc/systemd/system/monero.service
[Unit]
Description=monero full node
After=network-online.target

[Service]
User=monero
Group=monero
WorkingDirectory=/var/lib/monero
StateDirectory=/var/lib/monero
LogsDirectory=/var/log/monero
Type=simple
ExecStart=/usr/local/bin/monerod \
    --non-interactive \
    --data-dir /var/lib/monero \
    --block-sync-size 200 \
    --enforce-dns-checkpointing \
    --igd disabled \
    --p2p-bind-ip 0.0.0.0 \
    --p2p-bind-port 18080 \
    --p2p-external-port 18080 \
    --rpc-bind-ip 127.0.0.1 \
    --rpc-bind-port 18089 \
    --config-file /etc/monero/monerod.conf
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
