# manifest/resonance/zaphod/etc/containers/systemd/graylog.container
[Unit]
Description=Graylog Server (Podman Quadlet)
After=network-online.target mongo.service opensearch.service
Wants=network-online.target mongo.service opensearch.service

[Container]
Image=docker.io/graylog/graylog:6.3
ContainerName=graylog
EnvironmentFile=/etc/graylog/graylog.env
Pod=log-aggregator-pod
# Health
HealthCmd=curl -sf http://127.0.0.1:9000/api/system/lbstatus
HealthInterval=30s
HealthRetries=20
HealthTimeout=5s
# Security
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
