# manifest/resonance/zaphod/etc/containers/systemd/mongo.container
[Unit]
Description=MongoDB for Graylog (Podman Quadlet)
After=network-online.target graylog.network
Wants=network-online.target graylog.network

[Container]
Image=docker.io/library/mongo:7
ContainerName=mongo
Network=graylog
Volume=mongo-data:/data/db
# Keep it small; Graylog stores only metadata here
Environment=MONGO_INITDB_DATABASE=graylog
# Health
HealthCmd=mongosh --quiet --eval "db.runCommand({ ping: 1 })"
HealthInterval=30s
HealthRetries=5
HealthTimeout=5s
# Security
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
